# deploy/anyabot.service
[Unit]
Description=Anyabot Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/anyabot
ExecStart=/root/anyabot/venv/bin/python /root/anyabot/bot.py
Restart=always
RestartSec=5
StartLimitIntervalSec=0
User=root
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/anyabot.env

[Install]
WantedBy=multi-user.target
